# systemd/yo.service
[Unit]
Description=cow say yo

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/bin/yo.py

[Install]
WantedBy=multi-user.target
